# launchagents/com.coloradocareassist.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.coloradocareassist.telegram-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shulmeister/mac-mini-apps/careassist-unified/scripts/run-telegram-bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shulmeister/mac-mini-apps/careassist-unified</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/usr/local/bin:/opt/homebrew/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/shulmeister</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- KEEP ALIVE: Auto-restart if crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Restart throttling: wait 10 seconds between restart attempts -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/shulmeister/logs/telegram-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shulmeister/logs/telegram-bot-error.log</string>

    <!-- Process settings -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Soft resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>
